# boot.service
[Unit]
Description=test boot

[Service]
ExecStart=/home/pi/boot/steam.sh

[Install]
WantedBy=multi-user.target
